# infra/preview/polar-preview-infra.path
[Unit]
Description=Watch for preview infra triggers

[Path]
DirectoryNotEmpty=/srv/preview-triggers

[Install]
WantedBy=multi-user.target
