# infra/preview/polar-preview-hibernate.timer
[Unit]
Description=Check for idle preview backends every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
