# infra/preview/polar-preview-hibernate.service
[Unit]
Description=Hibernate idle preview backends

[Service]
Type=oneshot
ExecStart=/srv/preview-tools/preview-hibernate.sh
