# infra/preview/polar-preview-infra.service
[Unit]
Description=Process preview infra triggers

[Service]
Type=oneshot
ExecStart=/srv/preview-tools/process-preview-triggers.sh
