# infra/preview/polar-preview-logs.service
[Unit]
Description=Polar Preview Log Viewer
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /srv/preview-tools/log-viewer.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
