# server/supervisord.worker_github.conf
[supervisord]
nodaemon=true

[program:worker_github]
command=arq run_worker.WorkerSettingsGitHubCrawl
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
